# files/system/usr/lib/systemd/user/user-flatpak-setup.timer.d/disable-live.conf
[Unit]
ConditionKernelCommandLine=!rd.live.image
